# assets/systemd/gc-chkr.service
[Unit]
Description=Grafana Cloud node-side healthcheck
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/gc-chkr/env
ExecStart=/usr/bin/gc-chkr check --quiet
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/gc-chkr /var/log/gc-chkr
