# assets/debian/prerm
#!/bin/sh
set -e

case "$1" in
  remove|deconfigure|purge)
    if command -v systemctl >/dev/null 2>&1; then
      systemctl disable --now gc-hc.timer >/dev/null 2>&1 || true
      systemctl stop gc-hc.service >/dev/null 2>&1 || true
    fi
    ;;
esac

exit 0
